# src/main/deb/debian/postinst
#!/bin/sh
# postinst script for rserver
#
# see: dh_installdeb(1)

set -e

. /usr/share/debconf/confmodule

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)

    	# Create rserver user if it doesn't exist.
      if ! id rserver > /dev/null 2>&1 ; then
        adduser --system --home /var/lib/rserver --disabled-password rserver
      fi

      # rserver file structure on Debian
      # /etc/rserver: configuration
      # /usr/share/rserver: executable
      # /var/lib/rserver: data runtime
      # /var/log: logs

      rm -f /usr/share/rserver
      ln -s /usr/share/rserver-admin-* /usr/share/rserver

      if [ ! -e /var/lib/rserver/conf ] ; then
        ln -s /etc/rserver /var/lib/rserver/conf
      fi

      rlibs=/var/lib/rserver/R
      mkdir -p $rlibs
      echo $rlibs
      rpkgs=`find $rlibs -maxdepth 3 -mindepth 3 -type d`
      echo $rpkgs
      for pkg in $rpkgs
      do
        echo "R package: " $pkg
        mkdir -p $rlibs/all
        pkg_name=`basename $pkg`
        echo "R package name: " $pkg_name
        if [ ! -e $rlibs/all/$pkg_name ] ; then
          echo "R package move: " $rlibs/all/$pkg_name
          mv $pkg $rlibs/all
        fi
      done
      mkdir -p /var/lib/rserver/R/all

      if [ ! -e /var/log/rserver ] ; then
        mkdir /var/log/rserver
        ln -s /var/log/rserver /var/lib/rserver/logs
      fi

      if [ ! -e /tmp/rserver ] ; then
        mkdir /tmp/rserver
      fi

      chown -R rserver:adm /var/lib/rserver /var/log/rserver /etc/rserver /tmp/rserver
      chmod -R 750 /var/lib/rserver /var/log/rserver /etc/rserver/ /tmp/rserver
      find /etc/rserver/ -type f | xargs chmod 640

      update-rc.d rserver defaults >/dev/null
      /etc/init.d/rserver start
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
